# systemd/health-check.timer
[Unit]
Description=Run health-check every 5 minutes

[Timer]
OnUnitActiveSec=300
Persistent=true

[Install]
WantedBy=timers.target
